# packaging/xlog2-syncd.service
# systemd *user* service for the xlog2 headless sync/backup peer.
#
# It runs as your own login user so it picks up the GUI's configuration and data
# directories automatically:
#   config   ~/.config/xlog2/layout.ini      (shared secret + QRZ credentials)
#   data     ~/.local/share/xlog2-syncd/     (its own replica + QRZ cache)
# To point it elsewhere, add args to ExecStart, e.g.
#   ExecStart=/usr/bin/xlog2-syncd --config %h/.config/xlog2/layout.ini --data-dir %h/.local/share/xlog2-syncd
#
# Enable it for your user:
#   systemctl --user daemon-reload
#   systemctl --user enable --now xlog2-syncd.service
#   journalctl --user -u xlog2-syncd -f      # follow its log
#
# To run WITHOUT a logged-in session (headless box, or so it survives logout and
# starts at boot), enable lingering for the user once — this starts the user
# manager at boot and brings up its enabled services with no login required:
#   sudo loginctl enable-linger <user>

[Unit]
Description=xlog2 headless logbook sync/backup peer
Documentation=https://github.com/yo6ssw/xlog2/blob/master/docs/logbook-sync.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/xlog2-syncd
# Restart on crashes, but not when it exits cleanly (e.g. bad config).
Restart=on-failure
RestartSec=5
# Logs go to the journal (the daemon writes to stdout/stderr).
SyslogIdentifier=xlog2-syncd

# Only hardening valid under an unprivileged *user* manager. The Protect*/
# Restrict* sandboxing options are deliberately omitted: they manipulate the
# capability bounding set / mount namespaces, which a user manager lacks the
# privilege for, and would abort the service with 218/CAPABILITIES.
NoNewPrivileges=true

[Install]
WantedBy=default.target
